# snowflake-notifier/snowflake-notifier.timer
[Unit]
Description=Daily 08:00 trigger for Snowflake notifier

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
